# escooter-tracker.service
[Unit]
Description=E-Scooter Data Collector
Documentation=https://github.com/your-repo/escooter-tracker
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=erike
WorkingDirectory=/home/erike/raspberry-pi-escooter-tracker
ExecStart=/home/erike/raspberry-pi-escooter-tracker/venv/bin/python3 /home/erike/raspberry-pi-escooter-tracker/collector.py

# Restart policy
Restart=on-failure
RestartSec=30

# Logging
StandardOutput=append:/home/erike/raspberry-pi-escooter-tracker/systemd.log
StandardError=append:/home/erike/raspberry-pi-escooter-tracker/systemd.log

# Resource limits (prevent runaway process)
MemoryMax=256M
CPUQuota=50%

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
